# autostart/imwwheel.desktop
[Desktop Entry]
Type=Application
Exec=imwheel
Hidden=false
X-MATE-Autostart-enabled=true
Name[ja_JP]=imwheel
Name=imwwheel
Comment[ja_JP]=
Comment=
